# systemd/ptt.service
[Unit]
Description=Push-to-Talk daemon
After=graphical-session.target

[Service]
ExecStart=%h/.local/bin/ptt
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
