# systemd/ptt.service
[Unit]
Description=Push-to-Talk daemon
After=graphical-session.target keyd.service

[Service]
ExecStart=/home/hoshi/Projects/ptt/target/release/ptt-daemon
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
